# vm-builder/etc/init.d/dangerzone
#!/sbin/openrc-run
name="Dangerzone init script"
start() {
    # Hostname
	echo "dangerzone" > /etc/hostname
	echo "127.0.0.1 dangerzone" >> /etc/hosts
	hostname dangerzone

	# Networking	
	cat > /etc/network/interfaces << EOF
auto lo
iface lo inet loopback

auto eth0
iface eth0 inet dhcp
    hostname dangerzone
EOF
	setup-interfaces -a
	echo -e "\n" | setup-dns 4.4.4.4
	rc-service networking restart

	# Timezone
	setup-timezone -z UTC

    # Create user
    /usr/sbin/adduser -D -u 1001 user

	# Load the dangerzone container
	sudo -u user podman load -i /etc/dangerzone-converter.tar.gz

	# Allow podman containers to run
	echo "user:100000:65536" >> /etc/subuid
    echo "user:100000:65536" >> /etc/subgid

	# SSH reverse tunnel to host
	/etc/setup-ssh.py &
}